# roles/deploy_arrs/files/sonarr.service
[Service]
User=media-srv
Group=media-srv
